# src/modules/hotspot/filesystem/root/etc/systemd/system/ratos-autohotspot-deferred.timer
[Unit]
Description=RatOS PI5 — delayed autohotspot retry (Wi-Fi scan readiness)

[Timer]
OnBootSec=90s
AccuracySec=10s
Unit=ratos-autohotspot-deferred.service

[Install]
WantedBy=timers.target
