# src/modules/hotspot/filesystem/root/etc/systemd/system/ratos-nm-wlan-ap-guard.timer
[Unit]
Description=RatOS PI5 — periodic NM guard for autohotspot Wi-Fi AP

[Timer]
OnBootSec=25s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=ratos-nm-wlan-ap-guard.service

[Install]
WantedBy=timers.target
