# src/modules/hotspot/filesystem/root/etc/systemd/system/ratos-nm-wlan-ap-guard.service
[Unit]
Description=RatOS PI5 — Wi-Fi guard: keep AP unmanaged while up, recover AP when orphaned
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ratos-nm-wlan-ap-guard.sh
