# src/modules/hotspot/filesystem/root/etc/systemd/system/ratos-autohotspot-deferred.service
[Unit]
Description=RatOS PI5 — second autohotspot pass if hostapd did not stay up (boot race)
After=NetworkManager.service autohotspot.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ratos-autohotspot-deferred.sh
