# src/modules/hotspot/filesystem/root/etc/systemd/system/ravenos-nm-wlan-ap-guard.service
[Unit]
Description=RavenOS PI5 — keep Wi-Fi unmanaged while autohotspot hostapd AP is up
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ravenos-nm-wlan-ap-guard.sh
