# src/modules/hotspot/filesystem/root/etc/systemd/system/ravenos-autohotspot-deferred.timer
[Unit]
Description=RavenOS PI5 — delayed autohotspot retry (Wi-Fi scan readiness)

[Timer]
OnBootSec=90s
AccuracySec=10s
Unit=ravenos-autohotspot-deferred.service

[Install]
WantedBy=timers.target
