# src/modules/hotspot/filesystem/root/etc/systemd/system/autohotspot.service.d/10-after-nm.conf
[Unit]
# Bookworm: autohotspot uses nmcli + iw on wlan0 — start after NM is up.
After=NetworkManager.service

[Service]
# NM still probes wlan0 for a few seconds after boot; delay before hostapd reduces "AP then gone".
ExecStartPre=/bin/sleep 12
