# src/modules/mainsail/filesystem/root/etc/systemd/system/nginx.service.d/ratos-root.conf
[Service]
# Ensure /var/lib/ratos/nginx-mainsail-root.conf exists before nginx parses config.
ExecStartPre=/usr/local/bin/ratos-update-nginx-root.sh
