# src/modules/mainsail/filesystem/root/etc/systemd/system/nginx.service.d/ravenos-root.conf
[Service]
# Ensure /var/lib/ravenos/nginx-mainsail-root.conf exists before nginx parses config.
ExecStartPre=/usr/local/bin/ravenos-update-nginx-root.sh
